# supervisor/mongodb.conf
[supervisord]
nodaemon=false

[program:mongo]
directory=/data/db
command=/usr/bin/mongod --directoryperdb
autostart=true
autorestart=true
stopsignal=QUIT
